# plus/etc/profile.d/wine-tweaks.sh
export STAGING_SHARED_MEMORY=1
export STAGING_WRITECOPY=1
export STAGING_RT_PRIORITY_SERVER=90
